# copy/files/etc/init.d/mydhcpd
#!/sbin/openrc-run

command="/usr/sbin/dhcpd"
command_args="-q -pf /run/dhcpd.pid -user dhcp -group dhcp"  # Add any command line arguments here
pidfile="/run/dhcpd.pid"
name="DHCP Server"

depend() {
    after net.br-lan
}

